# web/treddy.service
[Unit]
Description=Treddy EMA Chart Service
After=network.target postgresql.service

[Service]
Type=simple
User=gabby
WorkingDirectory=/home/gabby/TradingBot/web
Environment="PGHOST=localhost"
Environment="PGPORT=5432"
Environment="PGDATABASE=trading_bot"
Environment="PGUSER=n8n_user"
Environment="PGPASSWORD=***REMOVED***"
ExecStart=/home/gabby/TradingBot/web/venv/bin/gunicorn --bind 0.0.0.0:5001 --workers 2 app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target